# workspace/crontab/higgses-payment-crond
* * * * * /usr/bin/php/var/www/testing/higgses/payment/current/artisan schedule:run >> /dev/null 2>&1