# workspace/crontab/zhangqiuzhe-testing-api-crond
* * * * * /usr/bin/php /var/www/testing/zhangqiuzhe/api/current/artisan schedule:run >> /dev/null 2>&1