# workspace/crontab/zhangqiuzhe-testing-seed-crond
* * * * * /usr/bin/php /var/www/testing/zhangqiuzhe/seed/current/artisan schedule:run >> /dev/null 2>&1